# HW2/Required Files/crontab
*/10 * * * * /./sorting /output_for_sort.txt 2>&1
